# systemd/dante-patchbox.service
[Unit]
Description=Dante Patchbox — AoIP matrix mixer and DSP patchbay
Documentation=https://github.com/legopc/dante-patchbox
After=network-online.target
Wants=network-online.target
# PTP clock daemon must be running before we start
After=statime.service ptp4l.service
Wants=statime.service

[Service]
Type=notify
User=patchbox
Group=patchbox

# Binary and config
ExecStart=/usr/local/bin/patchbox --config /etc/patchbox/config.toml

# Restart policy
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s

# Systemd watchdog — restart if process doesn't send WATCHDOG=1 within 10s
WatchdogSec=10s

# Environment
Environment=RUST_LOG=patchbox=info,tower_http=warn
Environment=INFERNO_USE_SAFE_CLOCK=false

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/patchbox /etc/patchbox

# Allow raw network sockets for Dante multicast
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW

# Real-time scheduling for the audio thread
LimitRTPRIO=99
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
